# pihole-stats.service
[Unit]
Description=pihole stats
After=syslog.target network.target auditd.service

[Service]
User=pi
ExecStartPre=/bin/sleep 30
ExecStart=/home/pi/piholedisplay/start-pihole-stat.sh

[Install]
WantedBy=multi-user.target
